# yabackup1/os/supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command = /bin/bash -c "/usr/sbin/crond -f -d 0"


[program:run_sh]
command=/bin/bash -c "/app/run.sh"


